# muschooly/muschooly.service
[Unit]
Description=Service 1
After=network.target

[Service]
User=mrblmg
WorkingDirectory=/home/mrblmg/sercives/fuzzy-octo-test/muschooly
ExecStart=/home/mrblmg/sercives/fuzzy-octo-test/muschooly/start.sh
Restart=always

[Install]
WantedBy=multi-user.target